# systemd/noaa-navionics.service
[Unit]
Description=Download NOAA ENC charts
StartLimitIntervalSec=6h
StartLimitBurst=3

[Service]
Type=oneshot
ExecStart=%h/.local/bin/noaa-navionics sync-charts --config %h/.config/noaa-navionics/config.ini --retries 5 --retry-delay 30
TimeoutStartSec=2h
Restart=on-failure
RestartSec=30min
