# systemd/noaa-navionics-preflight.service
[Unit]
Description=Write NOAA Navionics readiness report
Wants=noaa-navionics-track.service
After=noaa-navionics-track.service
StartLimitIntervalSec=30min
StartLimitBurst=60

[Service]
Type=oneshot
ExecStart=%h/.local/share/noaa-navionics/venv/bin/noaa-navionics status-report --config %h/.config/noaa-navionics/config.ini --gps-seconds-from-launcher-env %h/.config/noaa-navionics/launcher.env --output %h/.cache/noaa-navionics/status.json
TimeoutStartSec=15min
Restart=on-failure
RestartSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=default.target
