# systemd/noaa-navionics-track.service
[Unit]
Description=Log GPSD fixes to GPX with NOAA Navionics
StartLimitIntervalSec=10min
StartLimitBurst=60

[Service]
Type=simple
ExecStart=%h/.local/share/noaa-navionics/venv/bin/noaa-navionics log-track --config %h/.config/noaa-navionics/config.ini --rotate-daily
StandardOutput=null
StandardError=journal
Restart=on-failure
RestartSec=10
TimeoutStopSec=30s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=default.target
