# systemd/noaa-navionics-track.service
[Unit]
Description=Log GPSD fixes to GPX with NOAA Navionics

[Service]
Type=simple
ExecStart=%h/.local/bin/noaa-navionics log-track --config %h/.config/noaa-navionics/config.ini
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
