# systemd/noaa-navionics-preflight.service
[Unit]
Description=Write NOAA Navionics readiness report
StartLimitIntervalSec=5min
StartLimitBurst=5

[Service]
Type=oneshot
Environment=NOAA_NAVIONICS_GPS_SECONDS=10
EnvironmentFile=-%h/.config/noaa-navionics/launcher.env
ExecStart=%h/.local/bin/noaa-navionics status-report --config %h/.config/noaa-navionics/config.ini --gps-seconds ${NOAA_NAVIONICS_GPS_SECONDS} --output %h/.cache/noaa-navionics/status.json
TimeoutStartSec=0
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
