# systemd/noaa-navionics.service
[Unit]
Description=Download NOAA ENC charts
StartLimitIntervalSec=6h
StartLimitBurst=3

[Service]
Type=oneshot
ExecStartPre=%h/.local/bin/noaa-navionics wait-network --host www.charts.noaa.gov --port 443 --seconds 300
ExecStart=%h/.local/bin/noaa-navionics sync-charts --config %h/.config/noaa-navionics/config.ini --retries 5 --retry-delay 30
TimeoutStartSec=2h
Restart=on-failure
RestartSec=30min
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077
