# systemd/noaa-navionics-track.service
[Unit]
Description=Log GPSD fixes to GPX with NOAA Navionics
StartLimitIntervalSec=10min
StartLimitBurst=60

[Service]
Type=simple
ExecStart=%h/.local/bin/noaa-navionics log-track --config %h/.config/noaa-navionics/config.ini --rotate-daily
StandardOutput=null
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
