# systemd/noaa-navionics-preflight.service
[Unit]
Description=Write NOAA Navionics readiness report
StartLimitIntervalSec=5min
StartLimitBurst=5

[Service]
Type=oneshot
ExecStart=%h/.local/bin/noaa-navionics status-report --config %h/.config/noaa-navionics/config.ini --gps-seconds 10 --output %h/.cache/noaa-navionics/status.json
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
